# s6-overlay/s6-rc.d/godaddyupdater/run
#!/command/execlineb -S1
bash /etc/s6-overlay/scripts/run.sh